# ops/systemd/codex-usage-collector.timer
[Unit]
Description=Collect local AI session files for Codex Usage Profiler every five minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=5min
RandomizedDelaySec=60s
Persistent=true

[Install]
WantedBy=timers.target
